# pymultimonaprs.service
[Unit]
Description=PyMultimonAPRS
After=network.target

[Service]
ExecStart=/usr/bin/pymultimonaprs

[Install]
WantedBy=multi-user.target
