# pymultimonaprs.service
[Unit]
Description=PyMultimonAPRS
After=network.target

[Service]
ExecStart=pymultimonaprs

[Install]
WantedBy=multi-user.target
